# app-misc/kanata/files/kanata.initd
#!/sbin/openrc-run
# Copyright 1999-2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2 or later

supervisor=supervise-daemon
command="/usr/bin/kanata"
command_background="true"
